# crontab
0 */4 * * * /app/update.sh >> /var/log/cron.log 2>&1
